# server/scripts/debian/speedy-jobs.cron.d
0 4 * * * speedy /usr/lib/speedy-jobs/garbage_collector.py /etc/speedy-jobs.conf
5 * * * * speedy /usr/lib/speedy-jobs/maintain_tag_stats.py /etc/speedy-jobs.conf
