# server/speedy_notifier/debian/postrm
#!/bin/bash

update-rc.d speedy_notifier defaults
